# supervisord.conf
[supervisord]
nodaemon=true

[program:node]
command=yarn linter

[program:nginx]
command=nginx -g 'daemon off;'
