# supervisord.conf
[supervisord]
nodaemon=true

[program:node]
command=node dist/apps/gridiron/main.js

[program:nginx]
command=nginx -g 'daemon off;'